# supervisor/broker.ini
[program:kafka-broker]
command=/bin/su kafka /bin/bash -c "/opt/kafka/script/start.sh"
process_name=%(program_name)s
pidfile=/var/run/kafka-broker.pid
numprocs=1
autorestart=false
autostart=true
startsecs=10
startretries=1
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true